# templates/F_00_PRE_00_disable_apt_interactive_mode/etc/profile.d/zz99-disable_apt_interactive_mode.sh
# 1. setup needrestart.conf
# 2. setup ucf.conf
# And add the var, to setup once for all (especially for postfix)
export DEBIAN_FRONTEND="noninteractive"
