# templates/F_00_PRE_04_setup_login_msg/etc/update-motd.d/99-z-welcome
#!/bin/sh

COLOR_YELLOW="\033[1;33m"
COLOR_GREEN="\033[1;32m"
THIS_COLOR="${COLOR_GREEN}"
echo ""
echo "${THIS_COLOR}"
echo "=========================================================================================="
echo "                   You are now in \"$(hostname)\""
echo "=========================================================================================="
echo '\033[00m'
echo ""
